# modules/elexis/templates/gollum.init.erb
#!/bin/sh -e
### BEGIN INIT INFO
# Provides:          gollum
# Required-Start:    $local_fs  $remote_fs $network
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: Gollum: a nice simple Wiki
### END INIT INFO
#
# Start the gollum wiki

case "$1" in
  start) ;;
  stop|restart|force-reload) exit 0 ;;
  *) echo "Usage: $0 {start|stop|restart|force-reload}" >&2; exit 1 ;;
esac

<% inst_dir = scope.function_hiera(['elexis::praxis_wik::vcsRoot', '/home/elexis/praxis_wiki']) %>
if [ -d <%=  inst_dir %> ]
then 
  # redirect stdout and stderr into a log file
   sudo -iHu elexis gollum <%=  inst_dir %> &> /var/log/gollum.log &
fi
                                                                                                                                                 