# static-modules/cockpit/templates/cockpit.init.erb
#!/bin/bash -e
### BEGIN INIT INFO
# Provides:          cockpit
# Required-Start:    $local_fs  $remote_fs $network
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: Elexis-Cockpit simple Sinatra Web page for some administration work
### END INIT INFO
#
# Start the elexis-cockpit

case "$1" in
  start) ;;
  stop|restart|force-reload) exit 0 ;;
  *) echo "Usage: $0 {start|stop|restart|force-reload}" >&2; exit 1 ;;
esac

if [ -d <%=  @vcsRoot %> ]
then	
  # redirect stdout and stderr into a log file
  sudo -iHu elexis <%=  @useMock ? 'COCKPIT_CONFIG=mock_config.yaml' : '' %> <%=  @vcsRoot %>/elexis-cockpit.rb &> /var/log/elexis-cockpit.log &
fi
